# etc/profile.d/sh.local
#Add any required envvar overrides to this file, it is sourced from /etc/profile
export QT_QPA_PLATFORM=xcb
export OBS_USE_EGL=1
